# autostart/.config/autostart/xmodmap.desktop
[Desktop Entry]
Name=Xmodmap
Exec=/usr/bin/xmodmap .Xmodmap
X-GNOME-Autostart-enabled=true
Type=Application
